# supervisor/supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
nodaemon=true
logfile=/var/www/html/supervisor/logs/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=5MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

# [program:laravel-worker]
# process_name=%(program_name)s_%(process_num)02d
# priority=1
# command=php ./app/artisan queue:work 
# autostart=true
# autorestart=true
# stopasgroup=true
# killasgroup=true
# numprocs=3
# redirect_stderr=true
# stdout_logfile=/var/www/html/supervisor/logs/worker.log
# stopwaitsecs=3600
# stdout_logfile_maxbytes=5MB

[program:laravel-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php ./app/artisan queue:work 
autostart=true
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/html/supervisor/logs/worker.log


# [program:laravel-schedule-cron]
# process_name=%(program_name)s_%(process_num)02d
# command=php ./app/artisan schedule:run --verbose --no-interaction
# autostart=true
# autorestart=true
# user=root
# numprocs=3
# redirect_stderr=true
# priority=2
# stdout_logfile=/var/www/html/supervisor/logs/cron.log
# stopwaitsecs=3600
# stdout_logfile_maxbytes=5MB
